# etc/local.d/raspberry-pi3.stop
#/bin/bash
#
# Raspberry Pi3 : Disable specific hardware
#
/usr/bin/hciconfig hci0 down
exit 0

